# res/orbitbg.service
[Unit]
Description=Dynamic solar system X11 wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/orbitbg
