# res/orbitbg.service
[Unit]
Description=Dynamic solar system X11 wallpaper
Requires=xdg-desktop-autostart.target

[Service]
Type=oneshot
ExecStart=/usr/bin/orbitbg
